# calibre-web/calibre-web.container
[Unit]
Description=Calibre-Web Quadlet

[Container]
Image=lscr.io/linuxserver/calibre-web:latest
ContainerName=Calibre-Web-Quadlet
PublishPort=8083:8083
Environment=
Volume=/opt/podmnts/calibre-web:/config:Z
Volume=/opt/podmnts/calibre-web/library:/books:Z
User=root
Group=root
PodmanArgs=-e TZ=Europe/Copenhagen

[Service]
Restart=always
TimeoutStartSec=900

[Install]
WantedBy=multi-user.target default.target