# vpnstack/gluetun.container
[Unit]
Description=Gluetun Quadlet

[Container]
Image=qmcgaw/gluetun:latest
ContainerName=Gluetun-Quadlet
